# systemd/portmon.service
[Unit]
Description=PortMon - Port Traffic Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/portmon live
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
